# dockers/serf/supervisord-restart.conf
[program:restart-proc]
command=/restart-proc.sh
numprocs=1
autostart=true
autorestart=true
